# scripts/systemd/homelab_copilot_agent.service.sample
# Sample systemd unit for homelab_copilot_agent
# Drop this into /etc/systemd/system/homelab_copilot_agent.service and enable.

[Unit]
Description=Homelab Copilot Agent (sample)
After=network.target docker.service

[Service]
User=homelab
WorkingDirectory=/home/homelab/eddie-auto-dev/homelab_copilot_agent
Environment=PORT=8085
Environment=OLLAMA_HOST=http://192.168.15.2:11434
# Ensure the agent can resolve the Postgres container hostname and reach the API
Environment=API_BASE_URL=http://172.17.0.1:8503
# Recommended: use container hostname 'eddie-postgres' and attach to same Docker network
Environment=DATABASE_URL=postgresql://postgres:postgres@eddie-postgres:5432/postgres

# Run the container via helper that retrieves secrets from the Secrets Agent
# Ensure SECRETS_AGENT_URL and SECRETS_AGENT_API_KEY are set in a drop-in or environment
# Clear previous ExecStart so systemd will accept the replacement in a drop-in
ExecStart=
ExecStart=/bin/bash -c '/home/homelab/eddie-auto-dev/estou-aqui/scripts/secrets-agent/run-homelab-copilot.sh'
Restart=always

[Install]
WantedBy=multi-user.target
